# deploy/com.koriel.telegram-translator-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.koriel.telegram-translator-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/koriel/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>src</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/koriel/Development/telegram-translator-bot</string>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/koriel/Library/Logs/telegram-translator-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/koriel/Library/Logs/telegram-translator-bot.log</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
